# root/etc/services.d/v2ray/run
#!/usr/bin/with-contenv bash

exec /usr/bin/v2ray -config /config/v2ray/config.json
